# Startup/timesync.conf
description "Syncs clock with desktop"
author "Westley Davis"

start on startup
exec /etc/init/mytime.sh
